# templates/supervisord.conf.j2
[supervisord]
nodaemon=true

{% for program in supervisor_programs %}
[program:{{ program.name }}]
command={{ program.command }}

{% endfor %}
